# TP4/backup.service
[UNIT]
Description=Daily backup script

[Service]
Type=simple
User=backup
WotkingDirectory=/srv
ExecStart=/srv/backup.sh